# systemd/ibc-gateway.service
[Unit]
Description=IBC-managed IB Gateway
After=network-online.target
Wants=network-online.target
ConditionPathExists=/srv/data/spx-spark/runtime/ibc/config.ini
StartLimitIntervalSec=3600
StartLimitBurst=60

[Service]
Type=simple
WorkingDirectory=/home/ubuntu/spx-spark
EnvironmentFile=-/home/ubuntu/spx-spark/.env
Environment=IBC_CONFIG=/srv/data/spx-spark/runtime/ibc/config.ini
Environment=IBC_LOG_DIR=/srv/data/spx-spark/logs/ibc
ExecStart=/home/ubuntu/spx-spark/scripts/start-ibc-gateway.sh
Restart=always
RestartSec=60
TimeoutStopSec=45
KillSignal=SIGTERM

[Install]
WantedBy=default.target
