# systemd/ibc-gateway.service
[Unit]
Description=IBC-managed IB Gateway
After=network-online.target
Wants=network-online.target
ConditionPathExists=/srv/data/spx-spark/runtime/ibc/config.ini
# Never give up on relogin. With ExistingSessionDetectedAction=secondary the
# IBC session yields to a manual phone/desktop login and exits; systemd must
# keep retrying (once per RestartSec) until the manual session is gone so the
# feed recovers by itself. A finite start limit here (previously 60/3600s,
# exactly one hour of 60s retries) made systemd permanently stop the service
# after a long manual session.
StartLimitIntervalSec=0

[Service]
Type=simple
WorkingDirectory=/home/ubuntu/spx-spark
EnvironmentFile=-/home/ubuntu/spx-spark/.env
Environment=IBC_CONFIG=/srv/data/spx-spark/runtime/ibc/config.ini
Environment=IBC_LOG_DIR=/srv/data/spx-spark/logs/ibc
ExecStart=/home/ubuntu/spx-spark/scripts/start-ibc-gateway.sh
Restart=always
RestartSec=60
TimeoutStopSec=45
KillSignal=SIGTERM

[Install]
WantedBy=default.target
