# systemd/spx-spark-order-map-status.timer
[Unit]
Description=Pre-open market status report every 30 minutes (code gates to Beijing 14:15 -> US open)

[Timer]
OnCalendar=Mon..Fri *-*-* 14..22:00,30:00 Asia/Shanghai
Persistent=false

[Install]
WantedBy=timers.target
