# systemd/spx-spark-order-map-status.timer
[Unit]
Description=Market status + order map combined push (every 15 min Beijing 07:30-13:30, every 30 min 14:00 -> next-day 02:00)

[Timer]
# Morning: Beijing 07:30-13:30 every 15 minutes
OnCalendar=Mon..Fri *-*-* 07:30,45:00 Asia/Shanghai
OnCalendar=Mon..Fri *-*-* 08..12:00,15,30,45:00 Asia/Shanghai
OnCalendar=Mon..Fri *-*-* 13:00,15,30:00 Asia/Shanghai
# Afternoon through overnight: every 30 minutes
OnCalendar=Mon..Fri *-*-* 14..23:00,30:00 Asia/Shanghai
OnCalendar=Tue..Sat *-*-* 00..01:00,30:00 Asia/Shanghai
Persistent=false

[Install]
WantedBy=timers.target
