# systemd/spx-spark-post-close-review.service
[Unit]
Description=SPX Spark post-close SPX/SPXW review
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/ubuntu/spx-spark
EnvironmentFile=-/home/ubuntu/spx-spark/.env
ExecStart=/home/ubuntu/spx-spark/scripts/run-post-close-review.sh --date auto --quiet-if-empty
