# systemd/spx-spark-openclaw-weixin-keepalive.timer
[Unit]
Description=Refresh OpenClaw Weixin contextToken for proactive alerts

[Timer]
OnBootSec=10min
# contextToken has been observed to expire in ~70 minutes; keep the refresh
# interval well below that so proactive pushes do not hit sendMessage ret=-2.
OnUnitActiveSec=30min
Unit=spx-spark-openclaw-weixin-keepalive.service
Persistent=true

[Install]
WantedBy=timers.target
