# systemd/spx-spark-maintenance-weekly.service
[Unit]
Description=SPX Spark weekly threshold-gated retention maintenance
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/ubuntu/spx-spark
EnvironmentFile=-/home/ubuntu/spx-spark/.env
Environment=PATH=/home/ubuntu/.local/bin:/home/ubuntu/.cargo/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/ubuntu/spx-spark/scripts/run-maintenance-weekly.sh
UMask=0077
Nice=15
IOSchedulingClass=idle
TimeoutStartSec=30min
