# systemd/spx-spark-post-close-review.timer
[Unit]
Description=Run SPX Spark post-close SPX/SPXW review after US close

[Timer]
OnCalendar=Mon..Fri *-*-* 17:15:00 America/New_York
Persistent=true

[Install]
WantedBy=timers.target
