# systemd/spx-spark-maintenance-weekly.timer
[Unit]
Description=Run SPX Spark weekly non-destructive retention audit

[Timer]
OnCalendar=Sun *-*-* 13:00:00 Asia/Shanghai
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
