# systemd/spx-spark-order-map-refresh.timer
[Unit]
Description=Check for material order-map level changes every 20 minutes (code gates to US open -> Beijing 23:30)

[Timer]
OnCalendar=Mon..Fri *-*-* 21..23:10,30,50:00 Asia/Shanghai
Persistent=false

[Install]
WantedBy=timers.target
