# systemd/spx-spark-order-map-refresh.timer
[Unit]
Description=Order map refresh every 30 minutes at :15/:45 (interleaved with status at :00/:30; code gates to Beijing 14:15 -> next-day 02:00)

[Timer]
OnCalendar=Mon..Fri *-*-* 14..23:15,45:00 Asia/Shanghai
OnCalendar=Tue..Sat *-*-* 00..01:15,45:00 Asia/Shanghai
Persistent=false

[Install]
WantedBy=timers.target
